# warema-bridge/rootfs/etc/services.d/warema-bridge/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Warema Bridge
# S6 service entrypoint (legacy services.d)
# ==============================================================================
set -e

exec /run.sh
